# root/etc/s6-overlay/s6-rc.d/svc-readsb/run
#!/usr/bin/with-contenv bash

set -euo pipefail

readsb_user=${READSB_USER:-root}
readsb_run_dir=${READSB_RUN_DIR:-/run/readsb}
tar1090_db_path=${TAR1090_DB_PATH:-/usr/local/share/tar1090/aircraft.csv.gz}
build_metadata_path=/usr/local/share/readsb/build-metadata.env

# Import build metadata captured from upstream readsb so these values are
# available to the service process in production.
if [[ -r "${build_metadata_path}" ]]; then
    set -a
    # shellcheck disable=SC1090
    . "${build_metadata_path}"
    set +a
fi

if [[ "${readsb_user}" != "root" ]]; then
    if [[ "${readsb_user}" == "abc" && -z ${PUID:-} ]]; then
        export PUID=911
    fi
    if [[ "${readsb_user}" == "abc" && -z ${PGID:-} ]]; then
        export PGID=911
    fi
fi

timestamp_utc() {
    date -u '+%Y-%m-%d %H:%M:%S UTC'
}

print_kv() {
    local key=$1
    local value=$2
    printf '[%s] %-15s %s\n' "$(timestamp_utc)" "${key}:" "${value}"
}

# Create runtime directory at service start in case one-shot init did not run yet.
if ! mkdir -p "${readsb_run_dir}"; then
    echo "ERROR: unable to create READSB_RUN_DIR (${readsb_run_dir})." >&2
    echo "If running as non-root, ensure /run tmpfs uses matching uid/gid." >&2
    exit 1
fi

readsb_args_source="default"
if [[ -n ${READSB_ARGS:-} ]]; then
    readsb_args_source="env"
    read -r -a readsb_args <<< "${READSB_ARGS}"
else
    readsb_args=(
        --net
        --write-json "${readsb_run_dir}"
        --write-json-every 1
        --db-file "${tar1090_db_path}"
    )
fi

has_readsb_opt() {
    local opt=$1
    local token

    for token in "${readsb_args[@]}"; do
        if [[ ${token} == "${opt}" || ${token} == "${opt}="* ]]; then
            return 0
        fi
    done

    return 1
}

# Keep key file-output defaults even when READSB_ARGS is user-specified.
if ! has_readsb_opt "--write-json"; then
    readsb_args+=(--write-json "${readsb_run_dir}")
fi

if ! has_readsb_opt "--db-file"; then
    readsb_args+=(--db-file "${tar1090_db_path}")
fi

# Resolve the value of --device-type from the effective args.
_dt_val=""
_dt_prev=""
for _dt_tok in "${readsb_args[@]}"; do
    if [[ ${_dt_prev} == "--device-type" ]]; then
        _dt_val=${_dt_tok}
        break
    fi
    if [[ ${_dt_tok} == "--device-type="* ]]; then
        _dt_val=${_dt_tok#--device-type=}
        break
    fi
    _dt_prev=${_dt_tok}
done
unset _dt_prev _dt_tok

# If device type is rtlsdr and no --gain was provided, default to --gain auto.
if [[ ${_dt_val} == "rtlsdr" ]] && ! has_readsb_opt "--gain"; then
    readsb_args+=(--gain auto)
fi

# --- Feed profile: append --net-connector args per active profile ---
# shellcheck source=/usr/local/bin/feed-profiles.sh
. /usr/local/bin/feed-profiles.sh
parse_feed_profiles

feed_active=false
if [[ ${#FEED_PROFILE_LIST[@]} -gt 0 ]]; then
    feed_active=true

    # Add a --net-connector for each active profile
    for _fp in "${FEED_PROFILE_LIST[@]}"; do
        _fc=$(get_feed_connector "${_fp}")
        if [[ -z "${_fc}" ]]; then
            echo "[svc-readsb] WARNING: Unknown feed profile '${_fp}' — skipping." >&2
            continue
        fi
        readsb_args+=(--net-connector "${_fc}")
    done

    # UAT input connector (only if FEED_UAT_INPUT is set, e.g. "dump978:30978")
    if [[ -n "${FEED_UAT_INPUT:-}" ]]; then
        _uat_host="${FEED_UAT_INPUT%%:*}"
        _uat_port="${FEED_UAT_INPUT##*:}"
        if [[ -n "${_uat_host}" ]] && [[ -n "${_uat_port}" ]]; then
            readsb_args+=(--net-connector "${_uat_host},${_uat_port},uat_in,silent_fail")
        fi
    fi

    # Shared feed tuning args — only add if user hasn't already specified them
    if ! has_readsb_opt "--net-beast-reduce-interval"; then
        readsb_args+=(--net-beast-reduce-interval 0.5)
    fi
    if ! has_readsb_opt "--net-heartbeat"; then
        readsb_args+=(--net-heartbeat 60)
    fi
    if ! has_readsb_opt "--net-ro-size"; then
        readsb_args+=(--net-ro-size 1280)
    fi
    if ! has_readsb_opt "--net-ro-interval"; then
        readsb_args+=(--net-ro-interval 0.2)
    fi
    if ! has_readsb_opt "--json-location-accuracy"; then
        readsb_args+=(--json-location-accuracy 2)
    fi
    if ! has_readsb_opt "--range-outline-hours"; then
        readsb_args+=(--range-outline-hours 24)
    fi
    if ! has_readsb_opt "--max-range"; then
        readsb_args+=(--max-range 450)
    fi
    if ! has_readsb_opt "--write-json-every"; then
        readsb_args+=(--write-json-every 1)
    fi

    # UUID file
    uuid_file="/config/feed-uuid"
    if [[ -f "${uuid_file}" ]]; then
        readsb_args+=(--uuid-file "${uuid_file}")
    fi

    # Lat/lon from FEED_LAT/FEED_LON (fallback if not already in READSB_ARGS)
    if [[ -n "${FEED_LAT:-}" ]] && ! has_readsb_opt "--lat"; then
        readsb_args+=(--lat "${FEED_LAT}")
    fi
    if [[ -n "${FEED_LON:-}" ]] && ! has_readsb_opt "--lon"; then
        readsb_args+=(--lon "${FEED_LON}")
    fi
fi

readsb_args_rendered=$(printf '%q ' "${readsb_args[@]}")
readsb_args_rendered=${readsb_args_rendered% }

build_version=${VERSION:-unknown}
if [[ -r /build_version ]]; then
    build_version=$(< /build_version)
fi

echo "[$(timestamp_utc)] readsb container startup configuration"
print_kv "build_version" "${build_version}"
print_kv "VERSION" "${VERSION:-unknown}"
print_kv "BUILD_DATE" "${BUILD_DATE:-unknown}"
print_kv "VCS_REF" "${VCS_REF:-unknown}"
print_kv "PUID" "${PUID:-unset}"
print_kv "PGID" "${PGID:-unset}"
print_kv "UMASK" "${UMASK:-unset}"
print_kv "TZ" "${TZ:-unset}"
print_kv "READSB_USER" "${readsb_user}"
print_kv "READSB_ARGS(raw)" "${READSB_ARGS:-unset}"
print_kv "READSB_ARGS(src)" "${readsb_args_source}"
print_kv "effective_PUID" "${PUID:-unset}"
print_kv "effective_PGID" "${PGID:-unset}"
print_kv "READSB_RUN_DIR" "${readsb_run_dir}"
print_kv "TAR1090_DB_PATH" "${tar1090_db_path}"
print_kv "LSIO_NON_ROOT_USER" "${LSIO_NON_ROOT_USER:-unset}"
print_kv "effective_ARGS" "${readsb_args_rendered}"
print_kv "FEED_PROFILES" "${FEED_PROFILES:-none}"
if [[ "${feed_active}" == "true" ]]; then
    print_kv "FEED_UUID" "$(cat /config/feed-uuid 2>/dev/null || echo 'missing')"
    print_kv "FEED_UAT_INPUT" "${FEED_UAT_INPUT:-disabled}"
    print_kv "FEED_STATS" "${FEED_STATS_ENABLED:-true}"
    for _fp in "${FEED_PROFILE_LIST[@]}"; do
        _fc=$(get_feed_connector "${_fp}")
        _ms=$(get_mlat_server "${_fp}")
        print_kv "FEED[${_fp}]" "${_fc}"
        if [[ -n "${_ms}" ]]; then
            print_kv "MLAT[${_fp}]" "${_ms} (use separate mlat-client container)"
        fi
    done
fi

if [[ ${_dt_val} == "rtlsdr" ]]; then
    usb_char_devs=0
    if [[ -d /dev/bus/usb ]]; then
        usb_char_devs=$(find /dev/bus/usb -type c 2>/dev/null | wc -l | tr -d '[:space:]')
    fi
    print_kv "USB_CHAR_DEVS" "${usb_char_devs}"

    if command -v lsusb >/dev/null 2>&1; then
        lsusb_compact=$(lsusb 2>/dev/null | paste -sd '; ' -)
        print_kv "LSUSB" "${lsusb_compact:-none}"
    else
        print_kv "LSUSB" "unavailable"
    fi

    if [[ ! -e /dev/bus/usb ]]; then
        echo "[$(timestamp_utc)] WARNING: /dev/bus/usb is missing. Local RTL-SDR access will fail unless using network-only input." >&2
    elif [[ ! -r /dev/bus/usb ]]; then
        echo "[$(timestamp_utc)] WARNING: /dev/bus/usb is not readable by the container. Check device mapping and runtime permissions." >&2
    elif [[ ${usb_char_devs} -eq 0 ]]; then
        echo "[$(timestamp_utc)] WARNING: /dev/bus/usb is present but contains no USB device nodes." >&2
        echo "[$(timestamp_utc)] WARNING: Check host USB passthrough support and that the SDR is physically attached." >&2
    fi
fi

if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
    lsiown -R "${readsb_user}:${readsb_user}" "${readsb_run_dir}"
    exec s6-setuidgid "${readsb_user}" /usr/local/bin/readsb "${readsb_args[@]}"
fi

if [[ ! -w ${readsb_run_dir} ]]; then
    echo "ERROR: READSB_RUN_DIR (${readsb_run_dir}) is not writable by current user." >&2
    echo "When forcing user mode, set /run tmpfs uid/gid to match that user." >&2
    exit 1
fi

exec /usr/local/bin/readsb "${readsb_args[@]}"
