# root/etc/s6-overlay/s6-rc.d/svc-readsb/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# ── Logging ───────────────────────────────────────────────────────────────────
SVC_NAME="svc-readsb"
# shellcheck source=/usr/local/bin/log-functions.sh
. /usr/local/bin/log-functions.sh

readsb_user=${READSB_USER:-abc}
readsb_run_dir=${READSB_RUN_DIR:-/run/readsb}
tar1090_db_path=${TAR1090_DB_PATH:-/usr/local/share/tar1090/aircraft.csv.gz}
build_metadata_path=/usr/local/share/readsb/build-metadata.env

# Import build metadata captured from upstream readsb so these values are
# available to the service process in production.
if [[ -r "${build_metadata_path}" ]]; then
    set -a
    # shellcheck disable=SC1090
    . "${build_metadata_path}"
    set +a
fi

# Create runtime directory at service start in case one-shot init did not run yet.
if ! mkdir -p "${readsb_run_dir}"; then
    log_fatal "unable to create READSB_RUN_DIR (${readsb_run_dir})."
    log_fatal "If running as non-root, ensure /run tmpfs uses matching uid/gid."
    exit 1
fi

readsb_args_source="default"
if [[ -n ${READSB_ARGS:-} ]]; then
    readsb_args_source="env"
    read -r -a readsb_args <<< "${READSB_ARGS}"
else
    readsb_args=(
        --net
        --write-json "${readsb_run_dir}"
        --write-json-every 1
        --db-file "${tar1090_db_path}"
    )
fi

has_readsb_opt() {
    local opt=$1
    local token

    for token in "${readsb_args[@]}"; do
        if [[ ${token} == "${opt}" || ${token} == "${opt}="* ]]; then
            return 0
        fi
    done

    return 1
}

# Keep key file-output defaults even when READSB_ARGS is user-specified.
if ! has_readsb_opt "--write-json"; then
    readsb_args+=(--write-json "${readsb_run_dir}")
fi

if ! has_readsb_opt "--db-file"; then
    readsb_args+=(--db-file "${tar1090_db_path}")
fi

# Resolve the value of --device-type from the effective args.
_dt_val=""
_dt_prev=""
for _dt_tok in "${readsb_args[@]}"; do
    if [[ ${_dt_prev} == "--device-type" ]]; then
        _dt_val=${_dt_tok}
        break
    fi
    if [[ ${_dt_tok} == "--device-type="* ]]; then
        _dt_val=${_dt_tok#--device-type=}
        break
    fi
    _dt_prev=${_dt_tok}
done
unset _dt_prev _dt_tok

# If device type is rtlsdr and no --gain was provided, default to --gain auto.
if [[ ${_dt_val} == "rtlsdr" ]] && ! has_readsb_opt "--gain"; then
    readsb_args+=(--gain auto)
fi

# Suppress per-message decoder output unless user explicitly wants it (--debug / no --quiet).
if ! has_readsb_opt "--quiet" && ! has_readsb_opt "--debug"; then
    readsb_args+=(--quiet)
fi

# Apply SDR device assignment (manual READSB_DEVICE or auto-detected SDR_1090_INDEX)
if [[ ${_dt_val} == "rtlsdr" ]] && ! has_readsb_opt "--device"; then
    if [[ -n "${READSB_DEVICE:-}" ]]; then
        readsb_args+=(--device "${READSB_DEVICE}")
        log_info "Using manually assigned 1090 MHz dongle: --device ${READSB_DEVICE}"
    elif [[ -n "${SDR_1090_SERIAL:-}" ]]; then
        # Prefer serial over index: serial is stable across reboots
        readsb_args+=(--device "${SDR_1090_SERIAL}")
        log_info "Using auto-detected 1090 MHz dongle: --device ${SDR_1090_SERIAL} (serial)"
    elif [[ -n "${SDR_1090_INDEX:-}" ]]; then
        readsb_args+=(--device "${SDR_1090_INDEX}")
        log_info "Using auto-detected 1090 MHz dongle: --device ${SDR_1090_INDEX} (index)"
    fi
fi

# --- Feed profile: append --net-connector args per active profile ---
# shellcheck source=/usr/local/bin/feed-profiles.sh
. /usr/local/bin/feed-profiles.sh
parse_feed_profiles

feed_active=false
if [[ ${#FEED_PROFILE_LIST[@]} -gt 0 ]]; then
    feed_active=true

    # Add a --net-connector for each active profile (with per-profile UUID)
    for _fp in "${FEED_PROFILE_LIST[@]}"; do
        _fc=$(get_feed_connector "${_fp}")
        if [[ -z "${_fc}" ]]; then
            log_warn "Unknown feed profile '${_fp}' -- skipping."
            continue
        fi
        # Inject per-profile UUID into the connector string
        _fp_uuid=$(get_profile_uuid "${_fp}")
        if [[ -n "${_fp_uuid}" ]]; then
            _fc="${_fc},uuid=${_fp_uuid}"
        fi
        readsb_args+=(--net-connector "${_fc}")

        # Add secondary/failover connector if the profile has one
        _fc2=$(get_feed_connector_secondary "${_fp}")
        if [[ -n "${_fc2}" ]]; then
            if [[ -n "${_fp_uuid}" ]]; then
                _fc2="${_fc2},uuid=${_fp_uuid}"
            fi
            readsb_args+=(--net-connector "${_fc2}")
        fi
    done

    # UAT input connector (only if FEED_UAT_INPUT is set, e.g. "dump978:30978")
    if [[ -n "${FEED_UAT_INPUT:-}" ]]; then
        _uat_host="${FEED_UAT_INPUT%%:*}"
        _uat_port="${FEED_UAT_INPUT##*:}"
        if [[ -n "${_uat_host}" ]] && [[ -n "${_uat_port}" ]]; then
            readsb_args+=(--net-connector "${_uat_host},${_uat_port},uat_in,silent_fail")
        fi
    fi

    # Shared feed tuning args -- only add if user hasn't already specified them
    if ! has_readsb_opt "--net-beast-reduce-interval"; then
        readsb_args+=(--net-beast-reduce-interval 0.5)
    fi
    if ! has_readsb_opt "--net-heartbeat"; then
        readsb_args+=(--net-heartbeat 60)
    fi
    if ! has_readsb_opt "--net-ro-size"; then
        readsb_args+=(--net-ro-size 1280)
    fi
    if ! has_readsb_opt "--net-ro-interval"; then
        readsb_args+=(--net-ro-interval 0.2)
    fi
    if ! has_readsb_opt "--json-location-accuracy"; then
        readsb_args+=(--json-location-accuracy 2)
    fi
    if ! has_readsb_opt "--range-outline-hours"; then
        readsb_args+=(--range-outline-hours 24)
    fi
    if ! has_readsb_opt "--max-range"; then
        readsb_args+=(--max-range 450)
    fi
    if ! has_readsb_opt "--write-json-every"; then
        readsb_args+=(--write-json-every 1)
    fi

    # Lat/lon from FEED_LAT/FEED_LON (fallback if not already in READSB_ARGS)
    if [[ -n "${FEED_LAT:-}" ]] && ! has_readsb_opt "--lat"; then
        readsb_args+=(--lat "${FEED_LAT}")
    fi
    if [[ -n "${FEED_LON:-}" ]] && ! has_readsb_opt "--lon"; then
        readsb_args+=(--lon "${FEED_LON}")
    fi
fi

readsb_args_rendered=$(printf '%q ' "${readsb_args[@]}")
readsb_args_rendered=${readsb_args_rendered% }

build_version=${VERSION:-unknown}
if [[ -r /build_version ]]; then
    build_version=$(< /build_version)
fi

log_info "readsb container startup configuration"
log_kv "build_version" "${build_version}"
log_kv "VERSION" "${VERSION:-unknown}"
log_kv "BUILD_DATE" "${BUILD_DATE:-unknown}"
log_kv "VCS_REF" "${VCS_REF:-unknown}"
log_kv "PUID" "${PUID:-unset}"
log_kv "PGID" "${PGID:-unset}"
log_kv "UMASK" "${UMASK:-unset}"
log_kv "TZ" "${TZ:-unset}"
log_kv "READSB_USER" "${readsb_user}"
log_kv "READSB_ARGS(raw)" "${READSB_ARGS:-unset}"
log_kv "READSB_ARGS(src)" "${readsb_args_source}"
log_kv "READSB_RUN_DIR" "${readsb_run_dir}"
log_kv "TAR1090_DB_PATH" "${tar1090_db_path}"
log_kv "LSIO_NON_ROOT_USER" "${LSIO_NON_ROOT_USER:-unset}"
log_kv "effective_ARGS" "${readsb_args_rendered}"
log_kv "FEED_PROFILES" "${FEED_PROFILES:-none}"
log_kv "FEED_LAT" "${FEED_LAT:-unset}"
log_kv "FEED_LON" "${FEED_LON:-unset}"
log_kv "AUTO_LOCATION" "${READSB_AUTO_LOCATION:-true}"
if [[ "${feed_active}" == "true" ]]; then
    log_kv "FEED_UAT_INPUT" "${FEED_UAT_INPUT:-disabled}"
    log_kv "FEED_STATS" "${FEED_STATS_ENABLED:-true}"
    for _fp in "${FEED_PROFILE_LIST[@]}"; do
        _fc=$(get_feed_connector "${_fp}")
        _fc2=$(get_feed_connector_secondary "${_fp}")
        _ms=$(get_mlat_server "${_fp}")
        _fp_uuid=$(get_profile_uuid "${_fp}")
        log_kv "FEED[${_fp}]" "${_fc}"
        if [[ -n "${_fc2}" ]]; then
            log_kv "FEED[${_fp}:secondary]" "${_fc2}"
        fi
        log_kv "UUID[${_fp}]" "${_fp_uuid:-missing}"
        if [[ -n "${_ms}" ]]; then
            log_kv "MLAT[${_fp}]" "${_ms} (use separate mlat-client container)"
        fi
    done
fi

if [[ ${_dt_val} == "rtlsdr" ]]; then
    usb_char_devs=0
    if [[ -d /dev/bus/usb ]]; then
        usb_char_devs=$(find /dev/bus/usb -type c 2>/dev/null | wc -l | tr -d '[:space:]')
    fi
    log_kv "USB_CHAR_DEVS" "${usb_char_devs}"

    if command -v lsusb >/dev/null 2>&1; then
        lsusb_compact=$(lsusb 2>/dev/null | paste -sd '; ' -)
        log_kv "LSUSB" "${lsusb_compact:-none}"
    else
        log_kv "LSUSB" "unavailable"
    fi

    if [[ ! -e /dev/bus/usb ]]; then
        log_warn "/dev/bus/usb is missing. Local RTL-SDR access will fail unless using network-only input."
    elif [[ ! -r /dev/bus/usb ]]; then
        log_warn "/dev/bus/usb is not readable by the container. Check device mapping and runtime permissions."
    elif [[ ${usb_char_devs} -eq 0 ]]; then
        log_warn "/dev/bus/usb is present but contains no USB device nodes."
        log_warn "Check host USB passthrough support and that the SDR is physically attached."
    fi
fi

# ── Exec readsb with syslog-format log prefixing ─────────────────────────────────
# Pipe readsb stdout/stderr through awk to tag every line with a syslog-format
# prefix including a timestamp (readsb outputs decoder data to stderr).
LOG_PIPE=$(log_pipe_cmd decoder)

if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
    lsiown -R "${readsb_user}:${readsb_user}" "${readsb_run_dir}"
    exec s6-setuidgid "${readsb_user}" /bin/bash -c "
        /usr/local/bin/readsb $(printf '%q ' "${readsb_args[@]}") 2>&1 \
            | ${LOG_PIPE}
    "
fi

if [[ ! -w ${readsb_run_dir} ]]; then
    log_fatal "READSB_RUN_DIR (${readsb_run_dir}) is not writable by current user."
    log_fatal "When forcing user mode, set /run tmpfs uid/gid to match that user."
    exit 1
fi

exec /bin/bash -c "
    /usr/local/bin/readsb $(printf '%q ' "${readsb_args[@]}") 2>&1 \
        | ${LOG_PIPE}
"
