# root/etc/s6-overlay/s6-rc.d/init-readsb-config/run
#!/usr/bin/with-contenv bash

set -euo pipefail

# ── Logging ───────────────────────────────────────────────────────────────────
SVC_NAME="init-readsb-config"
# shellcheck source=/usr/local/bin/log-functions.sh
. /usr/local/bin/log-functions.sh

readsb_user=${READSB_USER:-abc}
readsb_run_dir=${READSB_RUN_DIR:-/run/readsb}
readsb_config_dir=/config/readsb
readsb_globe_history_dir=${readsb_config_dir}/globe_history

# ── Default PUID/PGID for abc user ───────────────────────────────────────────
# When running as the default abc user without explicit PUID/PGID, default to
# the LSIO base image UID/GID (911) so downstream scripts and log output have
# consistent values.
if [[ "${readsb_user}" != "root" ]]; then
    if [[ "${readsb_user}" == "abc" && -z ${PUID:-} ]]; then
        export PUID=911
    fi
    if [[ "${readsb_user}" == "abc" && -z ${PGID:-} ]]; then
        export PGID=911
    fi
fi

# Persist PUID/PGID to s6 container environment for downstream services
_s6_env="/run/s6/container_environment"
printf '%s' "${PUID:-}" > "${_s6_env}/PUID"
printf '%s' "${PGID:-}" > "${_s6_env}/PGID"

mkdir -p \
    "${readsb_config_dir}" \
    "${readsb_globe_history_dir}" \
    "${readsb_run_dir}"

if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
    lsiown -R "${readsb_user}:${readsb_user}" \
        "${readsb_config_dir}" \
        "${readsb_run_dir}"
fi

# --- Feed profile initialization ---
# shellcheck source=/usr/local/bin/feed-profiles.sh
. /usr/local/bin/feed-profiles.sh
parse_feed_profiles

# Persist the effective FEED_PROFILES value to s6 container environment.
# parse_feed_profiles defaults to 'adsbexchange' when FEED_PROFILES is unset.
if [[ ${#FEED_PROFILE_LIST[@]} -gt 0 ]]; then
    _s6_env="/run/s6/container_environment"
    printf '%s' "$(IFS=,; echo "${FEED_PROFILE_LIST[*]}")" > "${_s6_env}/FEED_PROFILES"
fi

if [[ ${#FEED_PROFILE_LIST[@]} -gt 0 ]]; then
    # ── Per-profile UUID generation ───────────────────────────────────────────
    # Each aggregator gets its own UUID stored in /config/feed-uuid-<profile>.
    # Priority: FEED_UUID_<PROFILE> env > existing per-profile file > auto-generate
    _uuid_regex='^[A-Fa-f0-9]{8}-[A-Fa-f0-9]{4}-[A-Fa-f0-9]{4}-[A-Fa-f0-9]{4}-[A-Fa-f0-9]{12}$'

    for _profile in "${FEED_PROFILE_LIST[@]}"; do
        _profile_uuid_file=$(get_profile_uuid_file "${_profile}")
        _env_suffix=$(get_profile_env_suffix "${_profile}")
        _env_var_name="FEED_UUID_${_env_suffix}"

        # Check for profile-specific env var override
        _env_uuid="${!_env_var_name:-}"

        if [[ -n "${_env_uuid}" ]]; then
            echo "${_env_uuid}" > "${_profile_uuid_file}"
            log_info "${_profile} UUID set from ${_env_var_name}: ${_env_uuid}"
        elif [[ ! -f "${_profile_uuid_file}" ]]; then
            _new_uuid=$(cat /proc/sys/kernel/random/uuid)
            echo "${_new_uuid}" > "${_profile_uuid_file}"
            log_info "${_profile} UUID auto-generated: ${_new_uuid}"
        fi

        # Validate UUID format
        _uuid_val=$(tr -d '[:space:]' < "${_profile_uuid_file}")
        echo "${_uuid_val}" > "${_profile_uuid_file}"
        if ! [[ ${_uuid_val} =~ ${_uuid_regex} ]]; then
            log_warn "${_profile} UUID invalid (${_uuid_val}). Regenerating."
            _new_uuid=$(cat /proc/sys/kernel/random/uuid)
            echo "${_new_uuid}" > "${_profile_uuid_file}"
            _uuid_val="${_new_uuid}"
        fi

        # Set ownership
        if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
            lsiown "${readsb_user}:${readsb_user}" "${_profile_uuid_file}"
        fi

        log_info "${_profile} UUID: ${_uuid_val}"
    done

    log_info "Active feed profiles: ${FEED_PROFILES:-adsbexchange}"

    # Log feed status URLs for each active profile
    for _profile in "${FEED_PROFILE_LIST[@]}"; do
        _status_url=$(get_feed_status_url "${_profile}")
        if [[ -n "${_status_url}" ]]; then
            log_info "Verify ${_profile} feed: ${_status_url}"
        fi
    done

    # Create stats temp directory
    mkdir -p /run/readsb-stats
    if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
        lsiown -R "${readsb_user}:${readsb_user}" /run/readsb-stats
    fi
fi

# ── Auto-location detection ──────────────────────────────────────────────────
# Detect receiver location via IP geolocation when FEED_LAT/FEED_LON not set.
# Explicit FEED_LAT/FEED_LON always take priority.
receiver_lat="${FEED_LAT:-}"
receiver_lon="${FEED_LON:-}"
auto_location="${READSB_AUTO_LOCATION:-true}"

if [[ -z "${receiver_lat}" || -z "${receiver_lon}" ]] && [[ "${auto_location}" =~ ^(true|yes|1)$ ]]; then
    # NOTE: ip-api.com free tier only supports HTTP -- HTTPS requires a paid plan.
    # This request sends only the container's public IP; no PII is transmitted.
    if geo_json=$(curl -sS -m 5 "http://ip-api.com/json/?fields=status,query,lat,lon,city,regionName,country,isp,timezone" 2>/dev/null); then
        geo_status=$(printf '%s' "${geo_json}" | grep -o '"status":"[^"]*"' | cut -d'"' -f4)
        if [[ "${geo_status}" == "success" ]]; then
            geo_ip=$(printf '%s' "${geo_json}" | grep -o '"query":"[^"]*"' | cut -d'"' -f4)
            receiver_lat=$(printf '%s' "${geo_json}" | grep -o '"lat":[0-9.eE+-]*' | cut -d: -f2)
            receiver_lon=$(printf '%s' "${geo_json}" | grep -o '"lon":[0-9.eE+-]*' | cut -d: -f2)
            geo_city=$(printf '%s' "${geo_json}" | grep -o '"city":"[^"]*"' | cut -d'"' -f4)
            geo_region=$(printf '%s' "${geo_json}" | grep -o '"regionName":"[^"]*"' | cut -d'"' -f4)
            geo_country=$(printf '%s' "${geo_json}" | grep -o '"country":"[^"]*"' | cut -d'"' -f4)
            geo_isp=$(printf '%s' "${geo_json}" | grep -o '"isp":"[^"]*"' | cut -d'"' -f4)
            geo_tz=$(printf '%s' "${geo_json}" | grep -o '"timezone":"[^"]*"' | cut -d'"' -f4)
            log_info "auto-detected location via IP ${geo_ip}: ${receiver_lat}, ${receiver_lon}"
            log_info "approximate location: ${geo_city}, ${geo_region}, ${geo_country} (tz: ${geo_tz}, isp: ${geo_isp})"
        else
            log_warn "IP geolocation failed (status: ${geo_status:-unknown})"
        fi
    else
        log_warn "IP geolocation request failed -- network may be unavailable"
    fi
fi

# Persist detected location in s6 environment for downstream services
if [[ -n "${receiver_lat}" && -n "${receiver_lon}" ]]; then
    _s6_env="/run/s6/container_environment"
    printf '%s' "${receiver_lat}" > "${_s6_env}/FEED_LAT"
    printf '%s' "${receiver_lon}" > "${_s6_env}/FEED_LON"
    log_info "receiver location: ${receiver_lat}, ${receiver_lon}"
else
    log_warn "no receiver location set (set FEED_LAT/FEED_LON or enable READSB_AUTO_LOCATION)"
fi

# ── SDR dongle auto-detection ─────────────────────────────────────────────────
# Enumerate RTL-SDR dongles and persist device assignments for svc-readsb.
# shellcheck source=/usr/local/bin/sdr-enumerate.sh
. /usr/local/bin/sdr-enumerate.sh

if enumerate_sdr_dongles; then
    _s6_env="/run/s6/container_environment"
    log_info "RTL-SDR dongles detected: ${SDR_DONGLE_COUNT} (mode: ${SDR_MODE})"

    if [[ -n "${SDR_1090_INDEX}" ]]; then
        log_info "1090 MHz dongle: device ${SDR_1090_INDEX} (serial: ${SDR_1090_SERIAL:-n/a}, mode: ${SDR_MODE})"
        printf '%s' "${SDR_1090_INDEX}" > "${_s6_env}/SDR_1090_INDEX"
        if [[ -n "${SDR_1090_SERIAL}" ]]; then
            printf '%s' "${SDR_1090_SERIAL}" > "${_s6_env}/SDR_1090_SERIAL"
        fi
    fi

    if [[ "${SDR_UAT_DETECTED}" == "true" ]]; then
        log_info "978 MHz UAT dongle detected (serial: ${SDR_UAT_SERIAL:-n/a}) -- not used by readsb"
        if [[ -n "${FEED_UAT_INPUT:-}" ]]; then
            log_info "FEED_UAT_INPUT=${FEED_UAT_INPUT} -- UAT data will be ingested from dump978 sidecar"
        else
            log_info "To decode UAT, run blackoutsecure/dump978 as a sidecar and set FEED_UAT_INPUT=dump978:30978"
            log_info "See: https://github.com/blackoutsecure/docker-dump978"
        fi
    fi

    if [[ "${SDR_MODE}" == "auto" && ${SDR_DONGLE_COUNT} -ge 2 && "${SDR_UAT_DETECTED}" != "true" ]]; then
        log_warn "Multiple RTL-SDR dongles found but none tagged as UAT (serial containing '978' or 'uat')"
        log_warn "Set 1090 dongle manually: READSB_DEVICE=<index|serial>"
        log_warn "Tag UAT dongle serial: rtl_eeprom -d <N> -s 00000978"
    fi

    if [[ "${SDR_MODE}" == "auto" && ${SDR_DONGLE_COUNT} -ge 2 && -z "${SDR_1090_INDEX}" ]]; then
        log_warn "No 1090 MHz dongle identified among ${SDR_DONGLE_COUNT} devices. readsb will use the first device."
        log_warn "Set READSB_DEVICE=<serial|index> to explicitly assign the 1090 dongle."
    fi
else
    log_warn "SDR enumeration failed (lsusb not found or no USB devices)"
fi

# ── RTL-SDR USB device permissions ────────────────────────────────────────────
# The abc user cannot access /dev/bus/usb devices (owned by root).
# Fix permissions on RTL-SDR device nodes so readsb can open them as non-root.
if [[ "${readsb_user}" != "root" ]] && command -v lsusb >/dev/null 2>&1; then
    while IFS= read -r devpath; do
        if [[ -c "${devpath}" ]]; then
            chmod 0666 "${devpath}"
            log_info "fixed USB permissions: ${devpath}"
        fi
    done < <(lsusb 2>/dev/null | grep -iE '0bda:(2832|2838)' | \
             sed -nE 's|^Bus ([0-9]+) Device ([0-9]+):.*|/dev/bus/usb/\1/\2|p')
fi
