# root/etc/s6-overlay/s6-rc.d/svc-feed-stats/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# ── Structured logging ───────────────────────────────────────────────────────────
SVC_NAME="svc-feed-stats"
# shellcheck source=/usr/local/bin/log-functions.sh
. /usr/local/bin/log-functions.sh
# shellcheck source=/usr/local/bin/feed-profiles.sh
. /usr/local/bin/feed-profiles.sh

feed_stats_enabled=${FEED_STATS_ENABLED:-true}

# Self-disable: sleep when stats are explicitly disabled
if [[ "${feed_stats_enabled,,}" != "true" ]]; then
    log_info "Stats disabled (FEED_STATS_ENABLED=${feed_stats_enabled}). Sleeping."
    exec sleep infinity
fi

# ── Configuration ─────────────────────────────────────────────────────────────
JSON_DIR="${READSB_RUN_DIR:-/run/readsb}"
JSON_FILE="${JSON_DIR}/aircraft.json"
AGE_COMPLAIN=30
AGE_INTERVAL=1800

# Console stats interval (default 2 minutes)
STATS_LOG_INTERVAL=${STATS_LOG_INTERVAL:-120}

# ADSBx stats upload interval (default 5 seconds, matching official ADSBx feeder)
ADSBX_UPLOAD_INTERVAL=${ADSBX_UPLOAD_INTERVAL:-5}
REMOTE_URL="https://adsbexchange.com/api/receive/"
MAX_CURL_TIME=10

# ── ADSBx stats upload setup ─────────────────────────────────────────────────
# When adsbexchange is an active feed profile, upload RSSI/stats data to ADSBx
# using the SAME UUID as the beast feed (--net-connector). Both pathways share
# a single identity so ADSBx sees one feeder, not two.
parse_feed_profiles
adsbx_upload=false
adsbx_uuid=""

if has_feed_profile "adsbexchange"; then
    # Prefer s6 env (set by init-readsb-config) over re-reading the file
    adsbx_uuid="${FEED_UUID_ADSBEXCHANGE:-$(get_profile_uuid "adsbexchange")}"
    if [[ -n "${adsbx_uuid}" ]]; then
        adsbx_upload=true
        log_info "ADSBx stats upload enabled (UUID: ${adsbx_uuid})"
        log_info "ADSBx stats URL: https://www.adsbexchange.com/api/feeders/?feed=${adsbx_uuid}"
    else
        log_warn "adsbexchange profile active but no UUID found — stats upload disabled"
    fi
fi

if [[ "${adsbx_upload}" == "true" ]]; then
    log_info "Starting stats service (upload=${ADSBX_UPLOAD_INTERVAL}s, console=${STATS_LOG_INTERVAL}s, JSON_DIR=${JSON_DIR})"
else
    log_info "Starting stats service (console=${STATS_LOG_INTERVAL}s, JSON_DIR=${JSON_DIR})"
fi

# ── Temp directory for upload staging ─────────────────────────────────────────
STATS_TMP_DIR="/run/readsb-stats"
if [[ "${adsbx_upload}" == "true" ]]; then
    mkdir -p "${STATS_TMP_DIR}"
fi

# ── Wait for JSON file ───────────────────────────────────────────────────────
stat_count=0
while [[ ${stat_count} -lt 12 ]]; do
    if [[ -f "${JSON_FILE}" ]]; then
        break
    fi
    stat_count=$((stat_count + 1))
    sleep 5
done

if [[ ! -f "${JSON_FILE}" ]]; then
    log_error "${JSON_FILE} not found. Is readsb running?"
    exit 1
fi

json_stat=$(stat --printf="%Y" "${JSON_FILE}" 2>/dev/null || echo "0")
old_age=${AGE_COMPLAIN}
last_console_stats=0

# Random startup jitter to reduce load spikes on ADSBx servers
if [[ "${adsbx_upload}" == "true" ]]; then
    sleep "$(( RANDOM % ADSBX_UPLOAD_INTERVAL )).$(( RANDOM % 100 ))"
fi

# ── Main loop ─────────────────────────────────────────────────────────────────
# Loop at upload rate when uploading, otherwise at console stats rate
if [[ "${adsbx_upload}" == "true" ]]; then
    LOOP_INTERVAL=${ADSBX_UPLOAD_INTERVAL}
else
    LOOP_INTERVAL=${STATS_LOG_INTERVAL}
fi

while true; do
    sleep ${LOOP_INTERVAL} &
    wait $!

    now=$(date +%s)

    # ── JSON freshness check ──────────────────────────────────────────────
    new_stat=$(stat --printf="%Y" "${JSON_FILE}" 2>/dev/null || echo "")
    if [[ -n "${new_stat}" ]]; then
        json_stat=${new_stat}
    fi

    diff=$(( now - ${json_stat:-0} ))
    if [[ ${diff} -gt ${old_age} ]]; then
        log_warn "${JSON_FILE} not updated in ${diff}s. Is the decoder running?"
        old_age=$(( old_age + AGE_INTERVAL ))
    else
        old_age=${AGE_COMPLAIN}
        # readsb is producing fresh JSON — reset the crash-loop fail counter
        # so next restart (e.g. Watchtower update) gets zero backoff.
        _fc_file="/run/readsb-fail-count"
        if [[ -f "${_fc_file}" ]] && [[ "$(cat "${_fc_file}" 2>/dev/null)" != "0" ]]; then
            echo "0" > "${_fc_file}"
        fi
    fi

    # ── ADSBx stats upload ────────────────────────────────────────────────
    # Matches official ADSBx format: enriches aircraft.json with UUID and
    # RSSI stats, gzips, and POSTs to the ADSBx receive API.
    if [[ "${adsbx_upload}" == "true" ]]; then
        _tmp_in="${STATS_TMP_DIR}/tmp.json"
        _tmp_out="${STATS_TMP_DIR}/upload.json"
        _tmp_gz="${STATS_TMP_DIR}/upload.gz"

        # Atomic copy to avoid reading a partial write from readsb
        if cp "${JSON_FILE}" "${_tmp_in}" 2>/dev/null; then
            # Build enriched payload with UUID and RSSI stats
            if jq -c \
                --arg UUID "${adsbx_uuid}" \
                '
                  .uuid = $UUID
                | .rssi = ([ .aircraft[]?.rssi // empty ] | if length == 0 then 0 else add / length | floor end)
                | ."rssi-min" = ([ .aircraft[]?.rssi // empty ] | if length == 0 then 0 else min | floor end)
                | ."rssi-max" = ([ .aircraft[]?.rssi // empty ] | if length == 0 then 0 else max | floor end)
                ' < "${_tmp_in}" > "${_tmp_out}" 2>/dev/null
            then
                gzip -c < "${_tmp_out}" > "${_tmp_gz}" 2>/dev/null

                if curl -sS -m ${MAX_CURL_TIME} -X POST \
                    -H "adsbx-uuid: ${adsbx_uuid}" \
                    -H "Content-Encoding: gzip" \
                    --data-binary @"${_tmp_gz}" \
                    "${REMOTE_URL}" >/dev/null 2>&1
                then
                    log_debug "ADSBx stats upload OK"
                else
                    log_warn "ADSBx stats upload failed (curl exit $?)"
                fi
            fi
        fi

        rm -f "${_tmp_in}" "${_tmp_out}" "${_tmp_gz}"
    fi

    # ── Periodic console stats ────────────────────────────────────────────
    elapsed=$(( now - last_console_stats ))
    if [[ ${elapsed} -ge ${STATS_LOG_INTERVAL} ]]; then
        _ac_count=$(jq -r '.aircraft | length // 0' < "${JSON_FILE}" 2>/dev/null || echo "?")
        _ac_with_pos=$(jq -r '[.aircraft[] | select(.lat != null)] | length // 0' < "${JSON_FILE}" 2>/dev/null || echo "?")
        _msgs=$(jq -r '.messages // 0' < "${JSON_FILE}" 2>/dev/null || echo "?")
        log_info "stats: ${_ac_count} aircraft tracked (${_ac_with_pos} with position), ${_msgs} messages total"
        last_console_stats=${now}
    fi
done
